# backup/hosts/alpha/systemd/restic-check.timer
[Unit]
Description=weekly restic repository check

[Timer]
OnCalendar=Sun *-*-* 05:30:00
Persistent=true

[Install]
WantedBy=timers.target
